# docker/app/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm -F

[program:nginx]
command=nginx -g "daemon off;"

; Scheduler Laravel : équivalent conteneurisé de l'entrée cron
; "* * * * * php artisan schedule:run". schedule:work est un processus
; long-running qui déclenche le scheduler chaque minute (backups:run et
; toute autre tâche planifiée). Redémarré automatiquement s'il s'arrête.
[program:scheduler]
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/www/html/storage/logs/scheduler.log
stderr_logfile_maxbytes=5MB
